# supervisord_local.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=debug                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:celery_worker]
command=python -m celery -A swifthub worker -l info
directory=/usr/src/app
autostart=true
autorestart=true
startretries=5
startsecs=5


[program:projekti_ar_gunicorn]
command=gunicorn swifthub.wsgi:application -b 0.0.0.0:8000 --reload -w 2 --pid /tmp/gunicorn.pid
directory=/usr/src/app
autostart=true
autorestart=true
startretries=5
startsecs=0
numprocs=1